# docker/services/pulse/run
#!/usr/bin/env bash

cd /var/www/html
exec php artisan pulse:check
