# docker/services/horizon/finish
#!/usr/bin/env bash

echo "❌ Horizon stopped! ❌"
